# deploy/systemd/stackpilot-controller.service
[Unit]
Description=StackPilot Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stackpilot-controller
Group=stackpilot-controller
WorkingDirectory=/opt/stackpilot/current
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8787
Environment=STACKPILOT_PRODUCTION=1
Environment=STACKPILOT_COOKIE_SECURE=1
Environment=STACKPILOT_DATABASE_PATH=/var/lib/stackpilot-controller/stackpilot.sqlite3
Environment=STACKPILOT_MASTER_KEY_FILE=%d/master-key
Environment=STACKPILOT_AGENT_TLS_CERT_PATH=%d/agent-tls.crt
Environment=STACKPILOT_AGENT_TLS_KEY_PATH=%d/agent-tls.key
EnvironmentFile=/etc/stackpilot/controller.env
LoadCredential=master-key:/etc/stackpilot/master-key
LoadCredential=agent-tls.crt:/etc/stackpilot/agent-tls.crt
LoadCredential=agent-tls.key:/etc/stackpilot/agent-tls.key
ExecStart=/usr/bin/node --preserve-symlinks-main /opt/stackpilot/current/apps/controller/dist/server.js
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
StateDirectory=stackpilot-controller
StateDirectoryMode=0700
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
ProtectProc=invisible
ProcSubset=pid
RemoveIPC=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/stackpilot-controller

[Install]
WantedBy=multi-user.target
