# deploy/systemd/stackpilot-firewall-helper.socket
[Unit]
Description=StackPilot Controller-only firewall helper socket

[Socket]
ListenStream=/run/stackpilot-firewall-helper/helper.sock
SocketUser=root
SocketGroup=stackpilot-firewall
SocketMode=0660
DirectoryMode=0750
Accept=yes
MaxConnections=32
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
