# deploy/systemd/stackpilot-database-helper.service
[Unit]
Description=StackPilot database-helper
Requires=stackpilot-database-helper.socket
After=local-fs.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/stackpilot-database-helper/current
Environment=NODE_ENV=production
Environment=STACKPILOT_DATABASE_HELPER_STATE_DIR=/var/lib/stackpilot-database-helper
ExecStart=/usr/bin/node --preserve-symlinks-main /opt/stackpilot-database-helper/current/apps/database-helper/dist/main.js
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
StateDirectory=stackpilot-database-helper
StateDirectoryMode=0700
UMask=0077
PrivateTmp=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
ProtectProc=invisible
RemoveIPC=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
